# scripts/motioncam.service
[Unit]
Description=MotionCam Node.js App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/motioncam
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/root/motioncam/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
